# servidor4/minikube/minikube.service
[Unit]
Description=Runs minikube service on start.

[Service]
ExecStart=minikube start
ExecStop=minikube stop
Type=oneshot
User=servidor4
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target